# boersenbot_dashboard.service
[Unit]
Description=Börsenbot Flask Dashboard
After=network.target

[Service]
Type=notify
User=veit
WorkingDirectory=/home/veit/boersenbot
ExecStart=/home/veit/boersenbot/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5000 --timeout 120 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
